# scripts/linguakit_api.service
[Unit]
After=network.service

[Service]
ExecStart=/home/folerpa/Linguakit-with-REST-API/run_linguakit_api.sh

[Install]
WantedBy=default.target
